# deploy/mohini.service
[Unit]
Description=Mohini Agent OS Daemon
Documentation=https://github.com/darshjme/mohini
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mohini
Group=mohini
ExecStart=/usr/local/bin/mohini start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment
EnvironmentFile=-/etc/mohini/env
WorkingDirectory=/var/lib/mohini

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/mohini
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=false
RestrictRealtime=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
